# olde/AST.service
[Unit]
Description=ROBOT
ConditionPathExists=|/bin
After=leds.service
 
[Service]
ExecStart=/bin/start
 
[Install]
WantedBy=multi-user.target